# .config/sv/mpdas/run
#!/bin/sh
exec /usr/local/bin/mpdas >/dev/null 2>&1
