# .config/sv/alsa-config-manager/run
#!/bin/sh
# exec $HOME/bin/alsa-config-manager
