# .config/sv/mpd/run
#!/bin/sh

# export ALSA_CONFIG_PATH=/usr/share/alsa/alsa.conf:/var/lib/bmix/bmix.conf
XDG_RUNTIME_DIR=/run/user/$(id -u) exec mpd --no-daemon ${XDG_CONFIG_HOME:-$HOME/.config}/mpd/mpd.conf
